# docker/c3b.conf
[program:c3b]
directory=/opt/cont3xt/c3b
command=npm start
autostart=true
autorestart=true
priority=100
