# docker/cont3xt.conf
[program:cont3xt]
directory=/opt/cont3xt
command=npm start
autostart=true
autorestart=true
priority=500
environment=PORT=3001
